# examples/quickstart/hello_agent/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3

[program:hello_agent]
command=python3 /app/hello_agent.py
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/hello_agent.err
stdout_logfile=/var/log/supervisor/hello_agent.out
